# scripts/ops/systemd/backend-refresh-scheduled.timer.example
[Unit]
Description=Run TerraPedia backend refresh on a schedule

[Timer]
OnBootSec=5min
OnUnitActiveSec=3h
Persistent=true
Unit=backend-refresh-scheduled.service

[Install]
WantedBy=timers.target
